# DenyHosts.cron
*/10 * * * * root python /usr/bin/denyhosts.py -c /etc/DenyHosts.cfg
